# packaging/brix-cache.service
# packaging/brix-cache.service
#
# Hardened systemd unit for the BriX-Cache storage gateway.
#
# This unit runs the system nginx binary with a dedicated BriX configuration
# file.  It is NOT a replacement for the stock nginx.service — it is a
# separate, independently managed instance for sites that want isolation.
#
# Installation:
#   install -m 644 packaging/brix-cache.service /etc/systemd/system/
#   # Create the dedicated config (copy from the example and edit paths/secrets):
#   cp /etc/nginx/conf.d/brix-cache.conf.example /etc/nginx/brix-cache.conf
#   # The config MUST contain:  pid /run/brix-cache.pid;
#   systemctl daemon-reload
#   systemctl enable --now brix-cache
#
# Logrotate note: the packaged rotate rule (/etc/logrotate.d/brix) references
# /run/nginx.pid (the stock instance PID file).  When using this unit, update
# that file to reference /run/brix-cache.pid instead, or copy it to a
# separate rotate rule for this instance.
#
# See: docs/09-developer-guide/deployment-hardening.md

[Unit]
Description=BriX-Cache storage gateway
Documentation=https://github.com/HEP-x/nginx-xrootd
After=network-online.target
Wants=network-online.target

[Service]
Type=forking
PIDFile=/run/brix-cache.pid
ExecStartPre=/usr/sbin/nginx -c /etc/nginx/brix-cache.conf -t
ExecStart=/usr/sbin/nginx -c /etc/nginx/brix-cache.conf
ExecReload=/bin/kill -s HUP $MAINPID
ExecStop=/bin/kill -s QUIT $MAINPID
# Align with brix-cache.conf.example recommendation (worker_shutdown_timeout 30s):
# an old worker finishes in-flight transfers then exits, and the native client
# reconnects within its stall window (XRDC_MAX_STALL_MS, default 30 s).
TimeoutStopSec=30

# File-creation umask: logs and staged uploads are not world-readable by default.
UMask=0027

# -------------------------------------------------------------------------
# Sandboxing — full rationale in docs/09-developer-guide/deployment-hardening.md
# -------------------------------------------------------------------------

# Prevents gaining new privileges through setuid/setgid executables or file
# capabilities at execve.  Does NOT block the broker's inline capability
# manipulation (capset + setresuid) — those operate on the permitted set
# already held, which NoNewPrivileges does not touch.
NoNewPrivileges=true

# /usr, /boot, /etc read-only.  nginx reads (never writes) its config and
# certificates from /etc; these paths never need write access.
ProtectSystem=strict

# Block access to /home, /root, /run/user.  If the export root lives under
# /home, add it to ReadWritePaths below instead.
ProtectHome=true

# Private /tmp namespace: staged upload partials use a tempfs-backed /tmp
# that is invisible to other processes and cleaned on service exit.
PrivateTmp=true

# No access to real device nodes.  nginx never opens /dev/sd*, /dev/tape, etc.
PrivateDevices=true

# Kernel tunable, module, and cgroup pseudo-filesystems are read-only.
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
ProtectKernelLogs=true

# Mount-namespace isolation (systemd builds a private mount namespace for the
# service before it drops privilege — the gateway itself cannot unshare(2) once
# its workers have shed CAP_SYS_ADMIN, so this is the correct layer for it).
#   ProtectProc=invisible  — the service sees only its OWN processes under /proc;
#                            other tenants' /proc/<pid> are hidden (a private,
#                            hidepid=invisible /proc mount in the namespace).
#   ProtectHostname / ProtectClock — the namespace cannot change the host name or
#                            the system clock.
#   ReadOnlyPaths          — the grid credential tree is read-only to the service
#                            (the host key is loaded once at startup; nothing needs
#                            to write it at runtime).
# NOTE: ProcSubset=pid (which also hides /proc/sys, /proc/meminfo, ...) is
# deliberately NOT set — some code paths read /proc/sys; enable it only after
# verifying on your build.
ProtectProc=invisible
ProtectHostname=true
ProtectClock=true
ReadOnlyPaths=/etc/grid-security

# Hide all kernel namespace operations.  The gateway does not create or enter
# namespaces; exposing unshare/setns/clone-with-ns only enlarges the attack surface.
RestrictNamespaces=true

# Prevent SCHED_FIFO / SCHED_RR acquisition.  All I/O is event-loop + thread-pool;
# no real-time scheduling is used or needed.
RestrictRealtime=true

# Prevent files created by the service from carrying setuid/setgid file bits.
# The impersonate broker changes its OWN process identity (setresuid/setresgid),
# which is governed by the @setuid syscall filter and CapabilityBoundingSet —
# not by this directive.  This correctly remains true (tight).
RestrictSUIDSGID=true

# Prevent mmap(PROT_WRITE|PROT_EXEC) and mprotect to PROT_EXEC.  nginx does not
# JIT-compile anything.  CAVEAT: if nginx was built with --with-pcre-jit,
# PCRE2's JIT engine requires W+X mappings and will crash with this directive
# enabled — in that case set MemoryDenyWriteExecute=false here.
MemoryDenyWriteExecute=true

# Prevent personality() calls that change the ABI of system calls; unnecessary
# for a native-Linux service and closes a class of seccomp-bypass gadgets.
LockPersonality=true

# Restrict socket families to the ones actually used:
#   AF_INET + AF_INET6  — root:// / WebDAV / S3 / metrics / TPC / CMS
#   AF_UNIX             — nginx worker-master IPC, thread-pool pipes
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX

# Allowlisted syscall groups:
#   @system-service  — baseline for any long-running daemon
#   @setuid          — setresuid / setresgid / setgroups for the impersonate broker
#
# Blocklisted groups (kernel-hardening):
#   @debug           — ptrace, perf_event_open; not needed at runtime
#   @mount           — mount/umount; not needed
#   @reboot          — reboot/kexec; definitely not needed
#   @swap            — swapon/swapoff; not needed
#   @obsolete        — old SysV/BSD calls; not needed
SystemCallFilter=@system-service @setuid
SystemCallFilter=~@debug @mount @reboot @swap @obsolete

# Block syscalls with a different ABI (e.g. 32-bit int 0x80 on x86_64).
SystemCallArchitectures=native

# -------------------------------------------------------------------------
# Capability bounding set
#
# Only three capabilities are retained in the bounding set:
#
#   CAP_SETUID, CAP_SETGID
#       Required by the impersonate broker.  The broker forks at startup
#       (lifecycle.c), drops its real uid/gid to a non-root service account
#       (broker_creds.c: imp_drop_to_service_user), and retains only these
#       two capabilities (imp_capset_setuid_setgid) so it can perform
#       per-request identity switches without running as root.  Workers
#       immediately call imp_worker_drop_caps() to shed these caps after
#       fork, so a compromised worker cannot escalate.
#
#       If xrootd_impersonation is disabled (or the module is built without
#       the impersonate component), these two caps can be removed:
#         CapabilityBoundingSet=CAP_NET_BIND_SERVICE
#       See docs/09-developer-guide/deployment-hardening.md §3.
#
#   CAP_NET_BIND_SERVICE
#       Allows the nginx master to bind ports < 1024.  XRootD listens on
#       1094 (< 1024); WebDAV on 443 / 8443; S3 on 9001 (> 1024 — only
#       CAP_NET_BIND_SERVICE is needed for 443/1094, 9001 does not need it).
#       If all listen ports are > 1024, this cap can also be dropped.
#
# CAP_DAC_OVERRIDE is explicitly NOT in the bounding set.  The broker's
# lifecycle.c lists it in kill_caps[] and drops it; nothing in the module
# uses it.
# -------------------------------------------------------------------------
CapabilityBoundingSet=CAP_SETUID CAP_SETGID CAP_NET_BIND_SERVICE

# -------------------------------------------------------------------------
# Writable filesystem paths
#
# Everything outside these paths is read-only (enforced by ProtectSystem=strict
# + ProtectHome=true + PrivateTmp=true).
#
# Adjust to the deployment's actual export root, stage dir, log dir, and
# any additional cert/keytab locations that nginx needs to write.
# -------------------------------------------------------------------------
#
#   /var/log/nginx       — access logs (xrootd_access*.log, WebDAV, S3)
#   /var/lib/brix-cache  — export root + default stage dir for partial uploads
#   /run                 — PID file (/run/brix-cache.pid), nginx lock files,
#                          and the impersonate broker UNIX socket
#
# Examples of additional paths to add:
#   /mnt/fastcache/xrd-stage   — if xrootd_stage_dir points here
#   /var/cache/brix-cache      — if nginx proxy_cache_path is configured
#   /home/xrootd               — if the export root is under /home (also
#                                requires ProtectHome=read-only or false)
ReadWritePaths=/var/log/nginx /var/lib/brix-cache /run

[Install]
WantedBy=multi-user.target
